# systemd/qa-app.service
[Unit]
Description=QA App Production Service
Documentation=https://github.com/qaapp/qaapp
After=network.target postgresql.service redis.service

[Service]
Type=forking
User=qaapp
WorkingDirectory=/opt/qa-app
Environment=NODE_ENV=production
EnvironmentFile=/opt/qa-app/.env
ExecStart=/usr/bin/pnpm run pm2:start
ExecReload=/usr/bin/pm2 reload ecosystem.config.js
ExecStop=/usr/bin/pnpm run pm2:stop
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qa-app

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/qa-app
PrivateTmp=true

[Install]
WantedBy=multi-user.target